# luci-app-campusportal/root/etc/uci-defaults/luci-i18n-campusportal-zh-cn
#!/bin/sh
uci -q batch <<-EOF >/dev/null
	set luci.languages.zh_cn='简体中文 (Simplified Chinese)'
	commit luci
EOF
rm -f /tmp/luci-indexcache /tmp/luci-modulecache/*
exit 0